# .config/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper rotation

[Service]
ExecStart=%h/.config/i3/scripts/wallpaper
Type=oneshot
