# .config/systemd/user/wallpaper.timer
[Unit]
Description=Wallpaper rotation timer

[Timer]
OnCalendar=7..20:0/5
Persistent=true

[Install]
WantedBy=default.target
